# CloudflareDDNS/rootfs/etc/services.d/ddns/run
#!/usr/bin/env bash
exec /usr/bin/cf_ddns
